# etc/cron.hourly/crx.syncMsQuotas
#!/bin/bash
# Copyright 2025 Peter Varkoly <pvarkoly@cephalix.et>

(
for user in $( /usr/sbin/crx_get_all_user.sh )
do
        echo -n "$user "
        /usr/sbin/crx_get_mquota.sh $user
done
) | gawk -f /usr/share/cranix/tools/quota.awk > /tmp/mquotalist
/usr/sbin/crx_api_post_file.sh users/syncMsQuotas /tmp/mquotalist
