# setup/systemd/kea-dhcp4.service.d/ExecStartPost.conf
[Service]
ExecStartPost=/usr/bin/sudo /usr/share/cranix/tools/create-room-classess.py
